# docker/supervisor/laravel-worker.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work
directory=/var/www/ARCVService/current_release
stdout_logfile=/var/www/ARCVService/storage/logs/worker.log
redirect_stderr=true
autostart=true
autorestart=true
user=apache
